# scripts/com.diego.bancotg.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.diego.bancotg.bot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-lc</string>
    <string>cd "PROJECT_DIR" &amp;&amp; source .venv/bin/activate &amp;&amp; python -m src.bot &gt;&gt; logs/bot.log 2&gt;&amp;1</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>PROJECT_DIR/logs/bot.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_DIR/logs/bot.stderr.log</string>
</dict>
</plist>
